# configs/systemd/cloudstack-lab-nat.service
[Unit]
Description=CloudStack lab NAT from cloudbr0 to Wi-Fi uplink
After=network-online.target NetworkManager-wait-online.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cloudstack-lab-nat.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
